# ops/systemd/project-mai-tai-prune-capture.service
[Unit]
Description=Prune central market-capture tables (market_capture_trades/quotes) older than 14 days
After=network-online.target

[Service]
Type=oneshot
User=trader
WorkingDirectory=/home/trader/project-mai-tai
EnvironmentFile=/etc/project-mai-tai/project-mai-tai.env
ExecStart=/home/trader/project-mai-tai/.venv/bin/python scripts/prune_market_ticks.py --keep-days 14 --tables market_capture_trades,market_capture_quotes
