# ops/systemd/project-mai-tai-trade-coach.service
[Unit]
Description=Project Mai Tai trade coach
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory=/home/trader/project-mai-tai
EnvironmentFile=/etc/project-mai-tai/project-mai-tai.env
Environment=MAI_TAI_TRADE_COACH_ENABLED=true
Environment=MAI_TAI_TRADE_COACH_REQUEST_TIMEOUT_SECONDS=20
Environment=MAI_TAI_TRADE_COACH_REVIEW_POLL_SECONDS=15
ExecStart=/home/trader/project-mai-tai/.venv/bin/mai-tai-trade-coach
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
UMask=0027
StandardOutput=append:/var/log/project-mai-tai/trade-coach.log
StandardError=append:/var/log/project-mai-tai/trade-coach.log
SyslogIdentifier=project-mai-tai-trade-coach

[Install]
WantedBy=multi-user.target
